# infra/systemd/attention.service.d/10-rt.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
